# wezterm-x/local.example/crontab
# Machine-local cron entries. Copy this file to wezterm-x/local/crontab,
# edit the schedule and message lines, then install with:
#
#   crontab wezterm-x/local/crontab
#
# Re-running `crontab <file>` replaces the user crontab atomically, so
# keeping the source of truth here means a fresh machine can restore the
# same reminders by rerunning that one command.
#
# Reminder entries delegate to scripts/runtime/reminder.sh, which finds
# the most recently active tmux client and pops a centered popup there
# (silently skips if no tmux client is attached). For raw popups beyond
# the reminder template, call scripts/runtime/tmux-popup-active.sh
# directly.
#
# Usage: reminder.sh <title> <message> [timeout-seconds] [width] [height]

# PATH: put user tmux (≥3.7 at ~/.local/bin) first, then /usr/local, then
# system. Cron's default /usr/bin/tmux is often 3.4 and cannot talk to a
# 3.7 server ("server exited unexpectedly" → silent popup no-op).
# See docs/tmux-install.md.
PATH=/home/REPLACE/.local/bin:/usr/local/bin:/usr/bin:/bin

# ⏰ 09:10 weekday punch-card reminder.
10 9 * * 1-5 /home/yuns/github/wezterm-config/scripts/runtime/reminder.sh ' ⏰ 打卡提醒 ' '打卡了吗？'

# Clean up stale managed-agent resume chains every 30 minutes. The script only
# targets long-running codex resume --last / claude --continue process groups
# with no controlling TTY by default, so attached tmux panes are left alone.
*/30 * * * * /home/yuns/github/wezterm-config/scripts/runtime/agent-cleanup.sh --kill --min-age 12h
